# stats/homepage-stats.service
[Unit]
Description=Homepage stats shim
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/services/web/homepage/stats
ExecStart=/usr/bin/node /srv/services/web/homepage/stats/app.mjs
Restart=always
RestartSec=3
User=apollo
SupplementaryGroups=docker

[Install]
WantedBy=multi-user.target
